# sources/mfa-sidecar-admin.service
[Unit]
Description=MFA Sidecar (Admin UI)
After=network-online.target mfa-sidecar-authelia.service
Wants=network-online.target

[Service]
Type=simple
User=mfa_sidecar
Group=mfa_sidecar
EnvironmentFile=-/etc/mfa-sidecar/mfa-sidecar.env
Environment=MFA_SIDECAR_POLICY_PATH=__INSTALL_DIR__/config/domain-policy.yaml
Environment=MFA_SIDECAR_RENDER_SCRIPT=__INSTALL_DIR__/bin/render_runtime_config.py
Environment=MFA_SIDECAR_STAGE_SCRIPT=__INSTALL_DIR__/bin/stage_runtime.py
Environment=MFA_SIDECAR_GENERATED_DIR=__INSTALL_DIR__/deploy/generated-runtime
Environment=MFA_SIDECAR_STAGE_ROOT=/
Environment=MFA_SIDECAR_ADMIN_BIND=127.0.0.1
Environment=MFA_SIDECAR_ADMIN_PORT=9087
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/usr/bin/python3 __INSTALL_DIR__/bin/admin_ui_app.py
Restart=on-failure
RestartSec=5s
PrivateTmp=true
ProtectHome=true
UMask=0027
WorkingDirectory=__INSTALL_DIR__

[Install]
WantedBy=multi-user.target
