# dist/package-root/sources/authelia.service
[Unit]
Description=MFA Sidecar (Authelia)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
EnvironmentFile=-/etc/mfa-sidecar/mfa-sidecar.env
ExecStart=/usr/local/bin/authelia --config /etc/mfa-sidecar/authelia/configuration.yml
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
